# src/sd-image/scripts/90-wireless.rules
ACTION=="add|change", SUBSYSTEM=="ieee80211", KERNEL=="phy0", RUN+="/usr/bin/iw phy %k interface add uap0 type __ap"
