# src/sd-image/scripts/electric-pi.service
[Unit]
Description=iCharger RESTful API Service
After=multi-user.target
Requires=multi-user.target

[Service]
Type=simple
ExecStart=/opt/run_electric.sh
ExecStop=/opt/stop_electric.sh
User=pirate
Restart=on-failure
RestartSec=8
WorkingDirectory=/opt

[Install]
WantedBy=multi-user.target
